# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=5MB
logfile_backups=2
pidfile=/tmp/supervisord.pid

[program:telegram_bot]
command=python3 /app/telegram_bot.py
directory=/app
autostart=true
autorestart=true
startretries=999
startsecs=5
stdout_logfile=/app/logs/telegram_bot.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_redirect=true
environment=PYTHONUNBUFFERED="1"

[program:scheduler]
command=python3 /app/agent.py --schedule --skip-checks
directory=/app
autostart=true
autorestart=true
startretries=999
startsecs=10
stdout_logfile=/app/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_redirect=true
environment=PYTHONUNBUFFERED="1"
